# systemd/pr-reviewer-bakeoff.timer
[Unit]
Description=Plow PR reviewer specialist bake-off (twice daily, 07:00 + 19:00 Pacific)

[Timer]
OnCalendar=*-*-* 07,19:00:00 America/Los_Angeles
Persistent=false
AccuracySec=1min
Unit=pr-reviewer-bakeoff.service

[Install]
WantedBy=timers.target
