# systemd/pr-reviewer-re-request.service
[Unit]
Description=Plow PR reviewer — translate GitHub re-request-review events into /review comments
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=odio
Group=odio
WorkingDirectory=/home/odio/.pr-reviewer
Environment=HOME=/home/odio
# PATH: system dirs FIRST (see pr-reviewer.service for rationale —
# ReadWritePaths grants write to ~/.local; without this order an
# attacker-placed ~/.local/bin/<command> would shadow system tools).
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/odio/.local/bin:/home/odio/.npm-global/bin
ExecStart=/home/odio/.pr-reviewer/re-request-poller.sh
TimeoutStartSec=5min
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pr-reviewer-re-request

NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHome=read-only
ReadWritePaths=/home/odio/.pr-reviewer /home/odio/.cache /home/odio/.local /home/odio/.config

[Install]
WantedBy=multi-user.target
