# systemd/pr-reviewer-bakeoff.timer
[Unit]
Description=Plow PR reviewer specialist bake-off (daily at 03:30 UTC)

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
AccuracySec=1min
Unit=pr-reviewer-bakeoff.service

[Install]
WantedBy=timers.target
