# systemd/pr-reviewer.service
[Unit]
Description=Plow PR reviewer (multi-specialist codex)
Documentation=https://github.com/cncorp/plow
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# KillMode=process so workers survive the orchestrator's exit. With the
# default `control-group` mode, systemd would kill every process in the
# unit's cgroup when ExecStart returns — including the detached workers
# review.sh just spawned. With `process`, only the main process (review.sh)
# is tracked for kill purposes; workers reparent to PID 1 and run to
# completion independently.
KillMode=process
User=odio
Group=odio
WorkingDirectory=/home/odio/.pr-reviewer
Environment=HOME=/home/odio
# PATH: system dirs FIRST so `timeout`, `bash`, `gh`, `git`, etc.
# resolve to /usr/bin paths regardless of any user-dir ordering.
# /home/odio/.local is NOT in ReadWritePaths (per the rationale on the
# ReadWritePaths line below), so ~/.local/bin tool replacements are
# blocked at the sandbox level too. The nvm-versioned bin dir provides
# `codex` (npm-managed under the operator's nvm default node version);
# bumping nvm node versions requires updating this path string in lock-
# step. User-installed tools (pipx, etc.) remain reachable via the
# trailing ~/.local/bin entry.
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/odio/.nvm/versions/node/v24.15.0/bin:/home/odio/.local/bin
# codex is npm-managed (see PATH above) and writes its cache on invocation.
# npm's default cache (~/.npm) sits under ProtectHome=read-only, so the write
# fails and codex stalls/errors. ~/.cache IS in ReadWritePaths (below), so
# point npm there. Keep this path inside a ReadWritePaths dir if you move it.
Environment=npm_config_cache=/home/odio/.cache/npm
ExecStart=/home/odio/.pr-reviewer/review.sh
TimeoutStartSec=90min
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pr-reviewer

# Memory ceiling for the whole unit cgroup (orchestrator + every
# detached worker, accumulated across ticks). Host is 128G; half is a
# comfortable upper bound that gets systemd-oomd involved before the
# unit starves the rest of the box. Observed pre-cap peak: 22.6G
# during the cncorp/plow chat-postgres cascade, so 60G/64G leaves
# >2x headroom for legitimate spikes while keeping runaway cascades
# self-terminating. MemoryHigh throttles first; MemoryMax is the hard
# stop where the kernel OOM-kills inside the cgroup.
MemoryHigh=60G
MemoryMax=64G

# Light sandbox — chosen to not break `just test` in tracked repos.
NoNewPrivileges=yes
# PrivateTmp deliberately omitted (load-bearing absence, not just "not set").
# review.sh dispatches each PR worker with `&` and exits (KillMode=process,
# Type=oneshot), so the unit deactivates within seconds while `just test`
# keeps running for up to 30m. PrivateTmp=yes makes systemd unmount the
# unit-private /tmp on deactivation — detached workers then see /tmp as
# ENOENT, and any `just test` that hardcodes /tmp (mktemp templates,
# Path("/tmp/...")) fails with false-positive `🧪 Tests failed` reviews.
# TMPDIR is pinned to $STATE_DIR/tmp upstream so well-behaved code keeps
# working, but that pin cannot help code that bypasses TMPDIR.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHome=read-only
# @KID_INDEX_RW_PATHS@ is rendered at install time by install.sh from
# repos.conf's KID_PATHS values, suffixed with `/.keepitdry` and
# prefixed with `-` so missing index dirs don't fail unit start. The
# kid prior-art lookup opens chromadb at
# `<KID_PATH>/.keepitdry/chroma.sqlite3`, and chromadb's SQLite backend
# needs to write WAL/journal files even on read-only queries. Without
# write access, queries fail with "attempt to write a readonly
# database" under ProtectHome=read-only and every review silently
# degrades to kid-less. Narrower than the kid-refresh render
# (@KID_RW_PATHS@): review specialists run with Codex's inner sandbox
# disabled and don't need writes outside the index dir. Do not edit
# @KID_INDEX_RW_PATHS@ by hand here.
# /home/odio/.local is intentionally NOT writable: it contains
# ~/.local/bin/<tool> symlinks (kid → uv-managed binary) and
# ~/.local/share/uv/tools/<tool>/bin/<tool> binaries that are PATH-search
# targets. PR-controlled `just test` runs in this service; if ~/.local
# were writable, an attacker could plant ~/.local/bin/codex or
# ~/.local/bin/kid and have the next reviewer/refresh tick exec it.
# Codex state lives in ~/.codex (separately writable above); claude
# state writes are from interactive operator sessions, not from this
# service. Add a narrower subdir if a future reviewer dependency
# legitimately needs to write under ~/.local.
ReadWritePaths=/home/odio/.pr-reviewer /home/odio/.codex /home/odio/.cache /home/odio/.config /home/odio/.docker @KID_INDEX_RW_PATHS@

[Install]
WantedBy=multi-user.target
