# systemd/pr-reviewer-approve.service
[Unit]
Description=Plow PR reviewer — /srosro-approve handler (post APPROVED reviews on trusted requests)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=odio
Group=odio
WorkingDirectory=/home/odio/.pr-reviewer
Environment=HOME=/home/odio
# PATH: system dirs FIRST (see pr-reviewer.service for rationale).
# This unit's ReadWritePaths is narrower (~/.pr-reviewer only), but
# keeping the same PATH order across all services makes the security
# stance uniform and survives future ReadWritePaths widenings.
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/odio/.local/bin:/home/odio/.npm-global/bin
ExecStart=/home/odio/.pr-reviewer/approve-from-replies.sh
TimeoutStartSec=5min
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pr-reviewer-approve

NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHome=read-only
ReadWritePaths=/home/odio/.pr-reviewer

[Install]
WantedBy=multi-user.target
