# systemd/pr-reviewer-re-request.service
[Unit]
Description=Plow PR reviewer — translate GitHub re-request-review events into /review comments
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=odio
Group=odio
WorkingDirectory=/home/odio/.pr-reviewer
Environment=HOME=/home/odio
# PATH: system dirs FIRST (see pr-reviewer.service for rationale —
# defense-in-depth on top of the ~/.local-not-in-ReadWritePaths fix).
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/odio/.nvm/versions/node/v24.15.0/bin:/home/odio/.local/bin
ExecStart=/home/odio/.pr-reviewer/re-request-poller.sh
TimeoutStartSec=5min
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pr-reviewer-re-request

NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHome=read-only
# See pr-reviewer.service for the rationale on excluding /home/odio/.local.
ReadWritePaths=/home/odio/.pr-reviewer /home/odio/.cache /home/odio/.config

[Install]
WantedBy=multi-user.target
